# dist/systemd/system/keylime_agent.service
[Unit]
Description=The Keylime compute agent
StartLimitInterval=10s
StartLimitBurst=5
Requires=var-lib-keylime-secure.mount
After=var-lib-keylime-secure.mount
Wants=tpm2-abrmd.service
After=tpm2-abrmd.service
# If the service should start only when hardware TPMs are available, uncomment the below lines
#ConditionPathExistsGlob=/dev/tpm[0-9]*
#ConditionPathExistsGlob=/dev/tpmrm[0-9]*

[Service]
ExecStart=/usr/bin/keylime_agent
KillSignal=SIGINT
TimeoutSec=60s
Restart=on-failure
RestartSec=120s
Environment="RUST_LOG=keylime_agent=info"
# If using swtpm with tpm2-abrmd service, uncomment the line below to set TCTI
# variable on the service environment
#Environment="TCTI=tabrmd:"

[Install]
WantedBy=default.target
